# salt/redis_server/redis-server.conf
# redis-server - upstart job file

description "redis server"
author "Ox Cart <ox@getlantern.org>"

start on (started tlsproxy_server and started tlsproxy_client)
stop on (stopping tlsproxy_server or stopping tlsproxy_client)
respawn

setuid redis
limit nofile 128000 128000
exec /usr/bin/redis-server /etc/redis/redis.conf
